# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:chromium]
command=chromium --no-sandbox --disable-setuid-sandbox --remote-debugging-address=0.0.0.0 --remote-debugging-port=9222 --disable-dev-shm-usage --disable-gpu --no-first-run --disable-background-networking --disable-sync --disable-default-apps --hide-scrollbars --mute-audio
autorestart=true
stdout_logfile=/var/log/supervisor/chromium.log
stderr_logfile=/var/log/supervisor/chromium-error.log
priority=10

[program:handoff]
command=node src/server.js
directory=/app
autorestart=true
stdout_logfile=/var/log/supervisor/handoff.log
stderr_logfile=/var/log/supervisor/handoff-error.log
environment=PORT="80",CDP_HTTP_URL="http://127.0.0.1:9222",SESSION_SECRET="%(ENV_SESSION_SECRET)s",PUBLIC_BASE_URL="%(ENV_PUBLIC_BASE_URL)s"
priority=20
startsecs=5